# templates/memcached.service.erb
[Unit]
Description=memcached daemon with SASL
After=network.target

[Service]
ExecStart=/usr/bin/memcached -u root -p <%= @port %> -m <%= @memory %> -S

[Install]
WantedBy=multi-user.target
